# debian/hello-go.service
[Unit]
Description=Hello Go Service

[Service]
ExecStart=/usr/bin/hello-go
Restart=always
User=root

[Install]
WantedBy=multi-user.target
